# tests/rabbitmq/supervisord.conf
[supervisord]
nodaemon=true

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
user=rabbitmq
environment=HOME="/var/lib/rabbitmq",USER="rabbitmq"

[program:flower]
command=/usr/bin/flower --port=5555
user=daemon
